# launchd/com.trening.bot.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.trening.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{APP_SUPPORT}}/scripts/start-bot.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>SHELL</key>
        <string>/bin/bash</string>
    </dict>

    <!-- Start ved login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Re-sjekk hvert minutt at sesjonen fortsatt finnes. Hvis tmux/claude
         har krasjet eller nett var nede ved boot, kommer den raskt opp igjen. -->
    <key>StartInterval</key>
    <integer>60</integer>

    <!-- Throttle slik at launchd ikke hamrer for fort mellom forsøk. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>{{LOGS}}/bot.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{LOGS}}/bot.stderr.log</string>

    <!-- Ikke blokker andre prosesser -->
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
